# descartes/etc/supervisor/conf.d/descartes.conf
[program:descartes]
command=/root/start/descartes.sh
autorestart=false
redirect_stderr=true