# descartes/etc/supervisor/conf.d/postgresql.conf
[program:postgres]
user=postgres
command=/usr/bin/postmaster -D "/var/lib/postgres/data"
stopsignal=INT
autorestart=true
redirect_stderr=true