# descartes/etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server
autorestart=false
redirect_stderr=true