# base/etc/supervisor/conf.d/collectd.conf
[program:collectd]
command=/root/start/collectd.sh
autorestart=true
redirect_stderr=true