# base/etc/supervisor/conf.d/rsyslog.conf
[program:rsyslog]
command=/root/start/rsyslog.sh
autorestart=true
redirect_stderr=true