# tomcat/etc/supervisor/conf.d/tomcat.conf
[program:tomcat7]
directory=/root/start
command=/root/start/tomcat.sh
autorestart=true
redirect_stderr=true