# graphiti/etc/supervisor/conf.d/graphiti.conf
[program:graphiti]
command=/root/start/graphiti.sh
autorestart=false
redirect_stderr=true